# bridge/com.notarangelo.coding-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notarangelo.coding-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>/Users/jacknotarangelo/Documents/GitHub/jax-assistant/bridge/coding-bridge.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jacknotarangelo/Documents/GitHub/jax-assistant/bridge</string>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/coding-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/coding-bridge.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_TLS_REJECT_UNAUTHORIZED</key>
        <string>0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
